# cron/crontab
0  8  *  *  *   /usr/local/sbin/check_raid.sh > /dev/null 2>&1
